# polymaker-frontend.service
[Unit]
Description=Polymaker Frontend Server
Documentation=https://github.com/polymaker
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/admin/Polymaker
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=FRONTEND_PORT=3000

ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 /home/admin/Polymaker/serve_frontend.py
ExecStop=/bin/kill -SIGTERM $MAINPID

# Restart policy (conservative)
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/var/log/polymaker/frontend.log
StandardError=append:/var/log/polymaker/frontend.log

# Resource limits
LimitNOFILE=1024
MemoryMax=256M

[Install]
WantedBy=multi-user.target
